# matrix.service
[Unit]
After=multi-user.target

[Service]
ExecStart=/home/pi/matrix/matrix
WorkingDirectory=/home/pi/matrix

[Install]
WantedBy=default.target
